# meeting-bot.service
[Unit]
Description=Autonomous Meeting Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/b24
ExecStart=/usr/bin/python3 /root/b24/autonomous_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=TELEGRAM_BOT_TOKEN=7992998044:AAHUNEIBP9nqIC7fmpHqBKAPcoQM5ltrWuI
Environment=GEMINI_API_KEY=AIzaSyDQR42zm4pcRMkY9KzKvEmXm7hyR8UzxHI
Environment=BITRIX_WEBHOOK_URL=https://skill-to-lead.bitrix24.ru/rest/1403/cmf3ncejqif8ny31/

[Install]
WantedBy=multi-user.target
