# meeting-bot.service
[Unit]
Description=Meeting Bot - Telegram Bot for Meeting Recording
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/meeting-bot
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/opt/meeting-bot
ExecStart=/usr/bin/python3 /opt/meeting-bot/meeting-bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Переменные окружения
EnvironmentFile=/opt/meeting-bot/.env

# Ограничения ресурсов
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
