# meeting-bot-monitor.service
[Unit]
Description=Meeting Bot Monitor - Auto-restart and monitoring
After=network.target meeting-bot.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/meeting-bot
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/opt/meeting-bot
ExecStart=/usr/bin/python3 /opt/meeting-bot/monitor_bot.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Переменные окружения
EnvironmentFile=/opt/meeting-bot/.env

# Ограничения ресурсов
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
