# systemd/meeting-bot.service
[Unit]
Description=Meeting Bot Assistant
Documentation=https://github.com/goqorhopar/b24
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/root/b24
ExecStart=/root/b24/venv/bin/python /root/b24/main.py
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStartSec=30
TimeoutStopSec=30
User=root
Group=root
Restart=always
RestartSec=10

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meeting-bot

# Переменные окружения
Environment=PYTHONPATH=/root/b24
Environment=LOG_LEVEL=INFO
Environment=USE_POLLING=true
Environment=PORT=3000

# Безопасность
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
